# toolchain/buildin/base-lite/common/etc/inittab
#
# /etc/inittab
#

#----- system start -----

::sysinit:/bin/mount -o rw,nosuid,nodev,noexec -t proc proc /proc
::sysinit:/sbin/ifconfig lo 127.0.0.1
::sysinit:/bin/hostname -F /etc/hostname
::wait:/bin/ash /etc/rc.mini

#----- system halt -----

::shutdown:/bin/ash /etc/rc.local.halt
