# toolchain/buildin/base-lite/common/etc/rc.local
#!/bin/ash

source /etc/rc.local.conf

echo -n 'starting local stuff:'

if [ X"${force_usb_mount}" != X"NO" ] ; then
	echo -n ' force usb-mount'
	/sbin/usbautomounter mount
fi

echo -n ' vdr'
export ${runvdr_vars}
/bin/ash ${runvdr} &

# With this sleep VDR start faster
/bin/sleep 5

exec /bin/ash /etc/rc.local.usb
